# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
childlogdir=/var/log/traefik

[program:traefik]
command=traefik

[program:rsyslogd]
command=rsyslogd
